# launchd/dev.vega-memory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.vega-memory</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/johnmacmini/workspace/vega-memory/dist/scheduler/index.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/johnmacmini/workspace/vega-memory</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/johnmacmini/workspace/vega-memory/data/logs/scheduler-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johnmacmini/workspace/vega-memory/data/logs/scheduler-stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>VEGA_DB_PATH</key>
    <string>/Users/johnmacmini/workspace/vega-memory/data/memory.db</string>
  </dict>
</dict>
</plist>
